# roles/samba4/sysvol/templates/glusterfsmounts.service.j2
[Unit]
After=glusterd.service
Wants=glusterd.service